# container-files/etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
stderr_logfile = NONE
stdout_logfile = /var/log/secure
redirect_stderr = true
